# supervisord.conf
[program:dali2mqtt]
command=/usr/bin/python3 dali-mqtt-daemon.py --config /app/
