# supervisord.conf
[program:dali2mqtt]
command=/usr/bin/python3 dali-mqtt-daemon.py --config /config/ --dali-lamps=3 --mqtt-server=192.168.5.120
